# docker-builds/reference-app/crontab
0 0 * * * root fetch_participant_resources /reference-research-app/development.ini >> /var/log/cron.log
#
