# docker-builds/reference-app/supervisord.conf
[supervisord]
nodaemon=true

[program:research-app]
command=pserve development.ini
directory=/reference-research-app

[program:cron]
command=cron -f
